# lib/systemd/system/linuxmuster-adminui.service
[Unit]
Description=linuxmuster.net AdminUI
After=network.target

[Service]
Type=forking
PIDFile=/var/run/ajenti.pid
ExecStart=/opt/linuxmuster/bin/python3 /opt/linuxmuster/bin/ajenti-panel -d --stock-plugins --plugins /usr/lib/linuxmuster-adminui/plugins

[Install]
WantedBy=multi-user.target
