# debian/postinst
#! /bin/sh
set -e

# Hide root warning for pip
export PIP_ROOT_USER_ACTION=ignore

msg()
{
    message=$1
    echo
    # Bold and orange font for linuxmuster.net
    echo "\e[1m\e[38;5;214m$message\e[39m\e[0m"
    echo
}

create_ajenti_venv()
{
    msg "Setting up virtual env in /opt"
    /usr/bin/python3 -m venv /opt/ajenti --system-site-packages
}

migrate_to_venv()
{
    # Checking if Ajenti is already installed in a venv
    if [ -f /opt/ajenti/bin/activate ] ; then
      return
    fi

    /usr/bin/python3 -m pip uninstall -y aj ajenti-panel ajenti.plugin.ace ajenti.plugin.core ajenti.plugin.dashboard ajenti.plugin.filesystem ajenti.plugin.plugins ajenti.plugin.settings ajenti.plugin.session-list
    create_ajenti_venv
}

case "$1" in
    install|configure)
        ajcfg="/etc/ajenti/config.yml"
        wucfg="/etc/linuxmuster/adminui/config.yml"
        adminuietc="/usr/lib/linuxmuster-adminui/etc"
        update_scripts="$adminuietc/update_scripts"
        install_scripts="$adminuietc/install_scripts"
        config_templates="$adminuietc/config_templates"

        # circumvent missing .installed files in package. Can be removed later
        if [ -f /etc/linuxmuster/.secret/administrator ]; then
            touch $adminuietc/.installed
        fi

        if [ ! -f $adminuietc/.installed ];
        then
            msg "---------------------------------------------------------------
Installing linuxmuster-adminui
---------------------------------------------------------------"

            create_ajenti_venv
            . /opt/ajenti/bin/activate
            PYTHON3=/opt/ajenti/bin/python3

            # do first install
            $PYTHON3 -m pip install -U pip wheel #setuptools #distribute
            hash -r
            /usr/bin/python3 -m pip install -r $adminuietc/requirements.txt

            msg "Copy default config files for adminui"
            mkdir -p /etc/linuxmuster/adminui/
            
            # Temporary fix for gipc and gevent in Ubuntu 24.04
            $PYTHON3 -m pip install gipc gevent -U
            
            #
            if [ -f $ajcfg ];then
                msg "Backing up config.yml"
                cp -a $ajcfg $ajcfg.bak
            fi

            # Update 7.1.11, same scheme for all directories, prefer underscore to hyphen
            if [ -d "/etc/linuxmuster/adminui/email-templates" ] ; then
                mv /etc/linuxmuster/adminui/email-templates /etc/linuxmuster/adminui/email_templates
            fi
            # Email templates will always be overwritten
            cp -a $adminuietc/email_templates /etc/linuxmuster/adminui/

            msg "Writing linuxmuster default ajenti config $ajcfg"
            cp -rf $config_templates/ajenti_default.yml $ajcfg
            if [ ! -f $wucfg ]; then
                # cp -n would be silent, but it's important to inform the user
                cp $config_templates/adminui_default.yml  $wucfg
            else
                msg "$wucfg already exists ... not overwriting with the default config file."
            fi

            # Avoid empty smtp config file until Ajenti 2.2.1
            if  [ ! -f /etc/ajenti/smtp.yml ] ; then
              cat << EOF > /etc/ajenti/smtp.yml
smtp:
  password:
  port:
  server:
  user:
EOF
            fi

            # enable systemctl service
            msg "Configure systemctl.."
            systemctl daemon-reload
            systemctl enable linuxmuster-adminui
            systemctl restart linuxmuster-adminui

            # generate permissions file
            $PYTHON3 $update_scripts/merge-permissions.py

            touch $adminuietc/.installed
            ipaddress=$(hostname --ip-address)
            msg "---------------------------------------------------------------
linuxmuster-adminui is now installed but not initialised!

Please open http://$ipaddress in a browser to end the configuration
----------------------------------------------------------------"

        else
            # do update
            msg "---------------------------------------------------------------
Updating linuxmuster-adminui
---------------------------------------------------------------"

            msg "Stop linuxmuster-adminui"
            systemctl stop linuxmuster-adminui

            migrate_to_venv
            . /opt/ajenti/bin/activate
            PYTHON3=/opt/ajenti/bin/python3

            msg "Upgrade pip to latest version"
            $PYTHON3 -m pip install pip -U

            msg "Remove deprecated plugins"
            $PYTHON3 -m pip uninstall -y ajenti.plugin.auth-users

            # Always backup config file if things are going bad
            msg "Backing up and updating config.yml"
            cp -a $ajcfg $ajcfg.bak
            cp -a $wucfg $wucfg.bak

            # Add new email parameters if needed ( Ajenti 2.1.44 )
            $PYTHON3 $update_scripts/update_ajenti_config.py

            # Split custom fields to /etc/linuxmuster/sophomorix/default-school (multischool-7.1.11)
            $PYTHON3 $update_scripts/clean_webui_config.py

            # Eventually move existing holiday file to /etc/linuxmuster/sophomorix/default-school (multischool-7.1.11)
            if [ -f /etc/linuxmuster/holidays.yml ] ; then
                mv /etc/linuxmuster/holidays.yml /etc/linuxmuster/sophomorix/default-school/holidays.yml
            fi

            msg "Check requirements"
            # Update Ajenti 2.1.44 : remove deprecated socketio library to avoid conflicts
            $PYTHON3 -c 'from socketio import mixins' 2>/dev/null && $PYTHON3 -m pip uninstall -y gevent-socketio-hartwork

            # Requirements must always be checked
<<<<<<< HEAD
<<<<<<< HEAD
            /usr/bin/python3 -m pip install -r $adminuietc/requirements.txt
            
            msg "Set default sophomorix adminui rights"
=======
=======
>>>>>>> d2ea5010c (Temporary fix for new Ubuntu.)
            $PYTHON3 -m pip install -r $webuietc/requirements.txt

            # Temporary fix for gipc and gevent in Ubuntu 24.04
            $PYTHON3 -m pip install gipc gevent -U

            msg "Set default sophomorix webui rights"
>>>>>>> d2ea5010c (Temporary fix for new Ubuntu.)
            $PYTHON3 $update_scripts/merge-permissions.py
            /usr/sbin/sophomorix-ui
            
            # Recreate upload folder using the right permissions
            msg "Verifying upload folder"
            bash $install_scripts/create_adminuiUploadFolder.sh

            # Temporary folder for webdav downloads
            mkdir -p /srv/.webdav
            chmod 777 /srv/.webdav

            # Avoid empty smtp config file until Ajenti 2.2.1
            if  [ ! -f /etc/ajenti/smtp.yml ] ; then
              cat << EOF > /etc/ajenti/smtp.yml
smtp:
  password:
  port:
  server:
  user:
EOF
            fi

            # start adminui
            msg "Start linuxmuster-adminui"
            systemctl daemon-reload
            systemctl start linuxmuster-adminui
            msg "---------------------------------------------------------------
Updating linuxmuster-adminui ...done
---------------------------------------------------------------"
        fi
    ;;
    upgrade|abort-upgrade)
    ;;
    
    *)
        msg "postinst called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

#DEBHELPER#

exit 0
